# collector/dev/systemd/systemd-cloudctl.service
# pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service
# autogenerated by Podman 1.6.4
# Sun Jul  5 18:33:06 UTC 2020

[Unit]
Description=Podman pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service
Documentation=man:podman-generate-systemd(1)
Requires=container-6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951.service container-eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6.service container-fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9.service
Before=container-6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951.service container-eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6.service container-fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start f280a032bf47919190e67ed1478605f5db60c9418bbd1736e63165dde3658e27
ExecStop=/usr/bin/podman stop -t 10 f280a032bf47919190e67ed1478605f5db60c9418bbd1736e63165dde3658e27
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/f280a032bf47919190e67ed1478605f5db60c9418bbd1736e63165dde3658e27/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6.service
# autogenerated by Podman 1.6.4
# Sun Jul  5 18:33:06 UTC 2020

[Unit]
Description=Podman container-eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service
After=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6
ExecStop=/usr/bin/podman stop -t 10 eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/eb5c99049a47f73425bca23a09fb28a9286ab5eb5bf031375770ea3daf6b2ab6/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9.service
# autogenerated by Podman 1.6.4
# Sun Jul  5 18:33:06 UTC 2020

[Unit]
Description=Podman container-fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service
After=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9
ExecStop=/usr/bin/podman stop -t 10 fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/fe22b0d414daba9795fde68165a4445ddb5c623db91739b5f3d185b23cb155a9/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951.service
# autogenerated by Podman 1.6.4
# Sun Jul  5 18:33:06 UTC 2020

[Unit]
Description=Podman container-6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service
After=pod-92221d7cb35d795d361295ee4ac2508437c6251fe6ff4176fc85c42e7e6a7430.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start 6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951
ExecStop=/usr/bin/podman stop -t 10 6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/6afb51f15117ea81cb2a48a6e2ec547830dd7a4e96543f5dd38bd2c4b9475951/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
